# deploy/xiaoboserver/campuskards-frontend.service
[Unit]
Description=CampusKards Frontend (Next.js)
After=network.target

[Service]
Type=simple
User=ck
Group=ck
WorkingDirectory=/ssd/7a40/campuskards/frontend

Environment=NODE_ENV=production
Environment=NEXT_PUBLIC_API_URL=https://gapi.xiaobocloud.fun

ExecStart=/usr/bin/node /ssd/7a40/campuskards/frontend/node_modules/.bin/next start \
  --port 3000

ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true

[Install]
WantedBy=multi-user.target
